# nestor-signalr-bridge/rootfs/etc/services.d/nestor-signalr-bridge/run
#!/usr/bin/with-contenv bashio

exec node /app/index.js
